# waterwheel/files/configs/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:playwright-mcp]
command=/usr/local/bin/playwright-mcp
autostart=false
autorestart=true
priority=1
stdout_logfile=/agent/outputs/playwright.log
stderr_logfile=/agent/outputs/playwright.err

[program:email-mcp]
command=/usr/local/bin/email-mcp
autostart=false
autorestart=true
priority=2
stdout_logfile=/agent/outputs/email-mcp.log
stderr_logfile=/agent/outputs/email-mcp.err